# data/psaas-update.service
[Unit]
Description=PSaaS Web Application Update Service

[Service]
Type=oneshot
ExecStart=/bin/bash /data/update.sh